# fan-control.service
[Unit]
Description=Temperature-based fan speed controller
After=network.target

[Service]
Type=simple
ExecStart={TARGETDIR}/venv/bin/python3 -u {TARGETDIR}/fan_control.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
